# svc/etc/env/prod/supervisord.conf
[supervisorctl]

[inet_http_server]
port = 127.0.0.1:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
logfile=/var/log/supervisord/supervisord.log    ; supervisord log file
logfile_maxbytes=1MB                            ; maximum size of logfile before rotation
logfile_backups=5                               ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=false                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=ec2-user                                   ; default user
childlogdir=/var/log/magenta                    ; where child log files will live

; This is where you run individual Tornado instances.
; We run four; one per processor core.
; In development, we ran as many as four per core with no issues.
; If you're looking to minimize cpu load, run fewer processes.
; BTW, Tornado processes are single threaded.
; To take advantage of multiple cores, you'll need multiple processes.

[program:dynoris]
user=ec2-user
priority=10
autostart=true
autorestart=false
environment=HOME="/home/ec2-user",USER="ec2-user",ASPNETCORE_ENVIRONMENT="Production",AWS_SHARED_CREDENTIALS_FILE="/home/ec2-user/.aws/credentials",AWS_CONFIG_FILE="/home/ec2-user/.aws/config"
command=dotnet dynoris.dll
directory=/srv/magenta/dynoris/dynoris/bin/Debug/netcoreapp2.1/
stderr_logfile = /var/log/dynoris/stderr.log
stdout_logfile = /var/log/dynoris/stdout.log
stderr_logfile_backups=7
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=7
stdout_logfile_maxbytes=10MB

[program:fe]
priority=300
autostart=true
autorestart=false
command=python /srv/magenta/svc/server.py --port=8092 --config_path=etc/conf/purple
directory=/srv/magenta/svc/
stderr_logfile = /var/log/magenta/stderr-fe.log
stdout_logfile = /var/log/magenta/stdout-fe.log
stderr_logfile_backups=10
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_logfile_maxbytes=10MB

[program:pay]
priority=300
autostart=false
autorestart=false
command=python /srv/magenta/svc/payments.py --port=8093 --config_path=etc/conf/purple
directory=/srv/magenta/svc/
stderr_logfile = /var/log/magenta/stderr-pay.log
stdout_logfile = /var/log/magenta/stdout-pay.log
stderr_logfile_backups=10
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_logfile_maxbytes=10MB

[program:poller]
priority=200
autostart=false
autorestart=false
numprocs=1
numprocs_start=1
process_name=poll-g
command=python /srv/magenta/svc/service.py --service=poller --log_path=/var/log/magenta/ --config_path=etc/conf/poll-prod --name=poll-prod --redis_db=0
directory=/srv/magenta/svc/
stderr_logfile = /var/log/magenta/stderr-poll-g.log
stdout_logfile = /var/log/magenta/stdout-poll-g.log
stderr_logfile_backups=10
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_logfile_maxbytes=10MB

[program:queue]
priority=150
autostart=false
autorestart=false
numprocs=1
numprocs_start=1
process_name=queue
command=python /srv/magenta/svc/service.py --service=queue --log_path=/var/log/magenta/ --config_path=etc/conf/purple --name=queue --redis_db=0
directory=/srv/magenta/svc/
stderr_logfile = /var/log/magenta/stderr-queue.log
stdout_logfile = /var/log/magenta/stdout-queue.log
stderr_logfile_backups=5
stderr_logfile_maxbytes=5MB
stdout_logfile_backups=5
stdout_logfile_maxbytes=5MB

[program:misc]
priority=150
autostart=false
autorestart=false
numprocs=1
numprocs_start=1
process_name=misc
command=python /srv/magenta/svc/service.py --service=misc --log_path=/var/log/magenta/ --config_path=etc/conf/purple --name=misc --redis_db=0
directory=/srv/magenta/svc/
stderr_logfile = /var/log/magenta/stderr-misc.log
stdout_logfile = /var/log/magenta/stdout-misc.log
stderr_logfile_backups=5
stderr_logfile_maxbytes=5MB
stdout_logfile_backups=5
stdout_logfile_maxbytes=5MB

[program:publisher-fb]
priority=100
autostart=false
autorestart=false
command=python /srv/magenta/svc/service.py --service=publisher --log_path=/var/log/magenta/ --config_path=etc/conf/purple --name=publisher-fb --provider=facebook,twitter --redis_db=0
directory=/srv/magenta/svc/
stderr_logfile = /var/log/magenta/stderr-pub-fb.log
stdout_logfile = /var/log/magenta/stdout-pub-fb.log
stderr_logfile_backups=10
stderr_logfile_maxbytes=20MB
stdout_logfile_backups=10
stdout_logfile_maxbytes=20MB

[program:publisher-fl]
priority=100
autostart=false
autorestart=false
command=python /srv/magenta/svc/service.py --service=publisher --log_path=/var/log/magenta/ --config_path=etc/conf/purple --name=publisher-fl --provider=linkedin,tumblr,flickr,500px --redis_db=0
directory=/srv/magenta/svc/
stderr_logfile = /var/log/magenta/stderr-pub-fl.log
stdout_logfile = /var/log/magenta/stdout-pub-fl.log
stderr_logfile_backups=10
stderr_logfile_maxbytes=20MB
stdout_logfile_backups=10
stdout_logfile_maxbytes=20MB

[group:fe]
programs=dynoris,fe
[group:pub]
programs=misc,queue,publisher-fb,publisher-fl
[group:poll]
programs=poller
